# bot-dashboard/dashboard.service
[Unit]
Description=Bot Dashboard API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/dashboard
ExecStart=/opt/dashboard/venv/bin/uvicorn api:app --host 0.0.0.0 --port 8090
Restart=always
RestartSec=5
Environment=DASHBOARD_TOKEN=0wp8-7IemmG0VxSqFuKlnw9REvF_LsFFUFwMDyn4NQ8

[Install]
WantedBy=multi-user.target
