# polyphemus/polyphemus.service
[Unit]
Description=Polyphemus Polymarket Signal Bot
After=network.target

[Service]
Type=notify
WorkingDirectory=/opt/polyphemus
ExecStart=/opt/polyphemus/venv/bin/python -m polyphemus.main
Restart=always
RestartSec=10
WatchdogSec=120
Environment=PYTHONUNBUFFERED=1
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectHostname=true
ProtectClock=true
ProtectHome=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
